# config/rules/99-scx-rusty-sync.rules
SUBSYSTEM=="power_supply", ACTION=="change", RUN+="/opt/scx-rusty-sync/scx-rusty-sync.sh" 